# roles/cinder-storage/templates/etc/systemd/system/cinder-volume.service.d/env.conf
[Service]
Environment="CONFIG_FILE=/etc/cinder/cinder-volume.conf"
